# luci-app-eqosplus/root/etc/hotplug.d/iface/10-eqosplus
#!/bin/sh
[ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "lan" ] && {
    sleep 5
    /etc/init.d/eqosplus restart
}
